# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:xvfb]
command=Xvfb :1 -screen 0 1024x768x24 -nolisten tcp
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

[program:x11vnc]
command=/scripts/start-vnc.sh
priority=20
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log
environment=DISPLAY=":1"

[program:jts-patcher]
command=/scripts/patch-trusted-ips.sh
priority=25
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/jts-patcher.log
stderr_logfile=/var/log/supervisor/jts-patcher.log

[program:ibc]
command=/scripts/start-gateway.sh
priority=30
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/ibc.log
stderr_logfile=/var/log/supervisor/ibc.log
environment=DISPLAY=":1"
